# task_files/IMESharePointDictionary__Microsoft_Office_.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2011-08-10T15:09:55.4511949</Date>
    <Description>This task updates SharePoint dictionary of Microsoft IME.</Description>
    <URI>\Microsoft\Office\IMESharePointDictionary</URI>
  </RegistrationInfo>
  <Principals>
    <Principal id="Author">
      <GroupId>S-1-5-32-545</GroupId>
    </Principal>
  </Principals>
  <Settings>
    <DisallowStartIfOnBatteries>true</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <RunOnlyIfIdle>true</RunOnlyIfIdle>
    <IdleSettings>
      <StopOnIdleEnd>true</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine>
    <MaintenanceSettings>
      <Period>P1D</Period>
      <Deadline>P7D</Deadline>
    </MaintenanceSettings>
  </Settings>
  <Triggers />
  <Actions Context="Author">
    <Exec>
      <Command>"c:\Program Files (x86)\Common Files\Microsoft Shared\IME16\IMESharePointDictionary.exe"</Command>
      <Arguments>-updateall</Arguments>
    </Exec>
  </Actions>
</Task>
