# task_files/Sync Controller Application__Microsoft_UE-V_.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Source>UE-V installer</Source>
    <Author>$(@%systemRoot%\system32\Microsoft.Uev.AppAgent.dll,-114)</Author>
    <Description>$(@%systemRoot%\system32\Microsoft.Uev.AppAgent.dll,-110)</Description>
    <URI>\Microsoft\UE-V\Sync Controller Application</URI>
  </RegistrationInfo>
  <Principals>
    <Principal id="Users">
      <DisplayName>Users Group</DisplayName>
      <GroupId>S-1-5-32-545</GroupId>
    </Principal>
  </Principals>
  <Settings>
    <AllowHardTerminate>false</AllowHardTerminate>
    <DisallowStartIfOnBatteries>true</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <ExecutionTimeLimit>PT4H</ExecutionTimeLimit>
    <Hidden>true</Hidden>
    <MultipleInstancesPolicy>Parallel</MultipleInstancesPolicy>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <Duration>PT10M</Duration>
      <WaitTimeout>PT1H</WaitTimeout>
      <StopOnIdleEnd>true</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
  </Settings>
  <Triggers>
    <LogonTrigger id="Logon trigger">
      <ExecutionTimeLimit>PT4H</ExecutionTimeLimit>
      <Delay>PT30S</Delay>
      <Repetition>
        <Interval>PT2M</Interval>
      </Repetition>
    </LogonTrigger>
  </Triggers>
  <Actions Context="Users">
    <Exec id="Check Settings Storage Location State">
      <Command>%windir%\system32\Microsoft.Uev.SyncController.exe</Command>
      <Arguments>$(Arg0) $(Arg1) $(Arg2) $(Arg3) $(Arg4)</Arguments>
      <WorkingDirectory>%windir%\system32</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
