# task_files/ThousandEyes Endpoint Agent Guardian__.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Guardian process to ensure that te-user-agent is always running.</Description>
    <URI>\ThousandEyes Endpoint Agent Guardian</URI>
  </RegistrationInfo>
  <Principals>
    <Principal id="AuthenticatedUsers">
      <GroupId>S-1-5-11</GroupId>
    </Principal>
  </Principals>
  <Settings>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <MultipleInstancesPolicy>Parallel</MultipleInstancesPolicy>
    <RestartOnFailure>
      <Count>3</Count>
      <Interval>PT1M</Interval>
    </RestartOnFailure>
    <StartWhenAvailable>true</StartWhenAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
  </Settings>
  <Triggers>
    <LogonTrigger>
      <Repetition>
        <Interval>PT15M</Interval>
      </Repetition>
    </LogonTrigger>
    <RegistrationTrigger />
  </Triggers>
  <Actions Context="AuthenticatedUsers">
    <Exec>
      <Command>"C:\Program Files (x86)\ThousandEyes\Endpoint Agent\te-guardian.exe"</Command>
      <Arguments>C:\Program Files (x86)\ThousandEyes\Endpoint Agent\te-user-agent.exe</Arguments>
    </Exec>
  </Actions>
</Task>
